# docker/sample-vnc-app/03-app.conf
[program:intellij]
command=%(ENV_HOME)s/app.sh
autostart=true