# docker/novnc/config/supervisord/00-xrdp.conf
[program:xrdp]
command=xrdp --nodaemon --config /etc/xrdp/xrdp.ini
autorestart=true


