# docker/novnc/config/supervisord/05-terminator.conf
[program:terminator]
command=/home/novnc/terminator.sh %(ENV_DISPLAY_1)s
autorestart=true
