# docker/novnc/config/supervisord/06-guaca-server.conf
[program:guacd]
command=/home/novnc/guacd.sh
autorestart=true

