# docker/novnc/config/supervisord/07-tomcat.conf
[program:tomcat]
command=/home/novnc/tomcat/bin/catalina.sh run
autorestart=true